# api/crontab
*/30 * * * * * * ./manage.py batch_order
